# system_files/etc/systemd/system/45d-firmware-check.timer
[Unit]
Description=Daily 45Drives Firmware Version Check
Documentation=https://github.com/45Drives/cockpit-hardware

[Timer]
# Run once per day with a random delay up to 12 hours to avoid thundering herd
OnCalendar=daily
RandomizedDelaySec=12h
# Run immediately on first boot if we missed a scheduled run
Persistent=true
# Don't wake the system from suspend
WakeSystem=false

[Install]
WantedBy=timers.target
