# packaging/ubuntu-jammy/postinst
#!/bin/bash

if [ "$1" == "configure" ]; then
	systemctl daemon-reload
	systemctl enable --now 45d-fancontroller.service 2>/dev/null || true
fi
