# system_files/usr/lib/systemd/system/45d-firmware-check.service
[Unit]
Description=45Drives Firmware Version Check
Documentation=https://github.com/45Drives/cockpit-hardware
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# firmware-check returns 2 when outdated firmware is detected (not an error)
SuccessExitStatus=2
ExecStart=/usr/bin/python3 /usr/share/45drives/firmware/firmware-check
# Run as root (needed for dmidecode, ipmitool, storcli, etc.)
User=root
# Limit resource usage
Nice=19
IOSchedulingClass=idle
# Timeout after 5 minutes
TimeoutStartSec=300
