# system_files/usr/lib/systemd/system/45d-fancontroller.service
[Unit]
Description=45Drives Fan Controller Daemon
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
ExecStart=/usr/share/cockpit/45drives-fan-controller/scripts/45d-fancontrollerd
Restart=unless-stopped
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=45d-fancontroller

# Hardening
ProtectSystem=strict
ReadWritePaths=/sys/class/hwmon /sys/bus/i2c/devices
ProtectHome=yes
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
